# pi/systemd/obd-led.service
# obd-led.service
#
# Status LED daemon — drives the two KY-016 RGB LEDs from collector state.
# Runs as its own long-lived process, fully decoupled from the collector:
# it only reads systemd unit status, sysfs, and the SQLite DB (read-only).
#
#   Type=simple    — no readiness handshake; the loop starts immediately.
#   Restart=always — recover from any crash; RestartSec=10 avoids thrashing.
#
# SupplementaryGroups=gpio grants access to /dev/gpiochip0 for the lgpio
# backend even if the user's primary groups change.

[Unit]
Description=OBD Status LEDs — Project 300K
After=network.target

[Service]
Type=simple
User=balu
SupplementaryGroups=gpio
WorkingDirectory=/home/balu/project300k/pi
ExecStart=/home/balu/project300k/pi/venv/bin/python src/led_status.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
# Hardening — defense-in-depth. ReadWritePaths=/mnt/usb is required even though
# the daemon only reads the DB: SQLite WAL readers must open the -shm sidecar
# read-write to map shared memory, so a read-only mount would break DB access.
# GPIO via /dev/gpiochip0 is unaffected by ProtectSystem (it does not cover /dev)
# and is reached through SupplementaryGroups=gpio above.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/mnt/usb
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
