# pi/systemd/obd-collector.service
# obd-collector.service
#
# Runs the OBD data collector as a persistent systemd service.
# Starts after Bluetooth is available so rfcomm0 can be bound.
#
# Key settings:
#   Type=notify      — main.py sends sd_notify READY=1 once init is complete.
#                      systemd waits for this before declaring the service active.
#   Restart=always   — restarts on any exit (crash, exception, clean exit).
#   RestartSec=15    — waits 15s before restarting to avoid thrashing on
#                      persistent errors (e.g. OBD dongle not found).
#   StartLimitBurst=5 — if the service restarts 5 times within
#   StartLimitIntervalSec=60s, systemd stops trying and marks it failed.
#                      Prevents infinite crash loops from consuming resources.
#   WatchdogSec=60   — if main.py does not send a watchdog ping within 60s,
#                      systemd kills and restarts the service. Catches deadlocks
#                      and infinite blocks that don't raise an exception.

[Unit]
Description=OBD Collector — Project 300K
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=notify
User=pi
WorkingDirectory=/home/pi/project300k/pi
ExecStart=/home/pi/project300k/pi/venv/bin/python src/main.py
Restart=always
RestartSec=15
StartLimitBurst=5
StartLimitIntervalSec=60
WatchdogSec=60
# OBD protocol detection can take 60-120s per attempt if the OBDLink is slow
# to respond. Default TimeoutStartSec=90s would kill the service before READY=1
# is sent on the first boot. 300s allows for multiple retries before giving up.
TimeoutStartSec=300
# TimeoutStopSec must exceed the OBD connection timeout (30s) plus queue drain
# time (up to 15s) so systemd does not send SIGKILL before shutdown completes.
TimeoutStopSec=60
# PYTHONUNBUFFERED ensures stdout/stderr lines reach journald immediately
# rather than being held in Python's block buffer in non-TTY mode.
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
