# pi/systemd/obd-datasette.service
[Unit]
Description=Datasette read-only browser for OBD SQLite data
After=network.target

[Service]
Type=simple
User=balu
ExecStart=/home/balu/project300k/pi/venv/bin/datasette --host 0.0.0.0 --port 8001 --immutable /mnt/usb/data/obd.db
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
