# pi/systemd/obd-collector.service
[Unit]
Description=OBD Collector — Project 300K
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=notify
User=pi
WorkingDirectory=/home/pi/project300k/pi
ExecStart=/home/pi/project300k/pi/venv/bin/python src/main.py
Restart=always
RestartSec=15
StartLimitBurst=5
StartLimitIntervalSec=60
WatchdogSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
