# pi/systemd/obd-sync.service
# obd-sync.service
#
# One-shot service that runs the sync script once per invocation.
# Triggered by obd-sync.timer — not started directly.
#
# Type=oneshot means systemd waits for the script to exit before
# marking the service as done. The timer will not fire again until
# the previous sync run has fully completed.
#
# Requires network-online.target so the sync script is not invoked
# before the network interface is up. The sync script performs its
# own two-step connectivity check (wlan0 IP + Tailscale ping) as a
# second layer of validation before attempting any POST.

[Unit]
Description=OBD Sync — Project 300K
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/project300k/pi
ExecStart=/home/pi/project300k/pi/venv/bin/python src/sync.py
# Belt-and-suspenders kill if sync hangs permanently (e.g. requests stall
# beyond the 30s timeout). Prevents a stuck sync from blocking the next
# timer-triggered invocation indefinitely.
TimeoutStartSec=120
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
