# pi/systemd/obd-sync.timer
# obd-sync.timer
#
# Fires obd-sync.service on a repeating schedule.
#
# OnBootSec=5min — first sync fires 5 minutes after boot. This delay
#   gives the iPhone hotspot time to connect and Tailscale time to
#   establish the VPN tunnel before the first sync attempt.
#
# OnUnitActiveSec=5min — subsequent syncs fire every 5 minutes after
#   the previous sync SERVICE ACTIVATION (not completion). This means if a
#   sync takes 3 minutes, the next one fires 2 minutes later. Use
#   OnUnitInactiveSec=5min instead if you want 5 minutes between completions.
#   Current choice: OnUnitActiveSec — ensures syncs attempt at least every
#   5 minutes regardless of how long each sync takes.

[Unit]
Description=OBD Sync Timer — Project 300K

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=obd-sync.service

[Install]
WantedBy=timers.target
