# pi/systemd/rfcomm-connect.service
# rfcomm-connect.service
#
# Establishes and maintains the Bluetooth RFCOMM connection to the OBDLink MX+.
# Runs rfcomm connect which blocks for the lifetime of the connection — systemd
# treats this as a persistent service. If the connection drops, Restart=always
# reconnects automatically.
#
# obd-collector.service orders after this unit (After=) and pulls it in with a
# soft Wants= — NOT Requires= — so the collector still starts (and retries the
# OBD link itself) even if this binding fails on the first attempt at boot.

[Unit]
Description=RFCOMM connection to OBDLink MX+
After=bluetooth.target
# Hard dependency on the BT stack — rfcomm cannot bind without bluetoothd.
Requires=bluetooth.service

[Service]
Type=simple
# Unblock rfkill and bring up hci0 before connecting.
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStartPre=-/usr/bin/hciconfig hci0 up
# Give the OBDLink MX+ time to finish booting its ELM327 chip after power-on
# before rfcomm attempts to connect. Without this delay, rfcomm connects at the
# BT layer before the ELM327 is ready, causing OBD initialization to fail.
ExecStartPre=-/usr/bin/rfcomm release 0
ExecStartPre=/bin/sleep 5
# rfcomm connect blocks until disconnected — this is intentional.
ExecStart=/usr/bin/rfcomm connect 0 00:04:3E:8A:94:4C
# Release the device cleanly on stop so the next start can bind it.
ExecStopPost=/usr/bin/rfcomm release 0
Restart=always
RestartSec=10
# Give hci0 time to come up before attempting connect.
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
